# c/vsmd/start-vsmd.conf
# This service maintains a vsm daemon.

respawn
console none

start on (local-filesystems)
stop on [!12345]

script
exec /usr/bin/vsmctl start
end script
